# templates/default/service.centos.erb
#!/bin/sh
# pyapns   This is the init script for starting up the PyAPNS
#
# chkconfig: 2345 20 80
# description: Starts and stops the PyAPNS  daemon that processes
#        all requests
# pidfile: /var/run/pyapns.pid
# File: pyapns
#
# About:
# Author: Sergey Balbeko 

# Source function library
. /etc/rc.d/init.d/functions

PYAPNS_PATH=<%= @home %>

script_result=0

start(){
  PYAPNS_START="Starting PyAPNS service: "
  echo -n "$PYAPNS_START"
  
  pid=`ps aux | grep "pyapns.server.APNSServer" | grep "twistd" | awk '{print $2 }'`
  if [ $pid ]
  then
    echo
    echo -n "PyAPNS is already running. "
    echo_failure
    echo
    script_result=1
    return
  fi
  
  pid=''
  
  twistd --logfile=/var/log/pyapns.log --pidfile=/var/run/pyapns.pid -r epoll web --class=pyapns.server.APNSServer --port=<%= @port %>
  
  sleep 2
  

  if [ -f "/var/run/pyapns.pid" ]
  then
    success "$PYAPNS_START"
    echo
  else
    failure "$PYAPNS_START"
    echo
    script_result=1
  fi
}

stop(){
  PYAPNS_STOP="Stopping PyAPNS service: "
  echo -n "$PYAPNS_STOP"
  if [ -f "/var/run/pyapns.pid" ]
  then
    pid=`head /var/run/pyapns.pid`
  else
    pid=''
  fi
  if [ $pid ]
  then
    k=`kill -s 2 $pid`
    sleep 1
    if [ -f "/var/run/pyapns.pid" ]
    then
      rm -f /var/run/pyapns.pid
    fi
  echo_success
  else
  echo_failure 
  script_result=1
  fi
  echo
}

restart(){
  stop
  sleep 1
  start
}

status(){
  if [ -f "/var/run/pyapns.pid" ]
  then
    pid=`cat /var/run/pyapns.pid`
    out=`ps aux | grep $pid | grep 'APNSServer'`
    
    if [ "$out" ]
    then
      echo "PyAPNS Server (pid $pid) is running..."
    else
      echo "PyAPNS Server is stopped"
    fi
  else
    echo "PyAPNS Server is stopped"
  fi
}

# See how we were called.
case "$1" in
  start)
  start
  ;;
  stop)
  stop
  ;;
  restart)
  restart
  ;;
  status)
  status
  ;;
  *)
  echo $"Usage: $0 {start|stop|restart|status}"
  exit 1
esac

exit $script_result
